# proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=crond -f -l 6 -L /proc/1/fd/1
redirect_stderr=true

[program:nginx]
command=nginx -g 'daemon off;'
redirect_stderr=true